# filechooser-portal/io.elementary.files.xdg-desktop-portal.service.in
[Unit]
Description=Portal service (Pantheon Files)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.jeremypw.files
ExecStart=@libexecdir@/com.github.jeremypw.dogfood-files-7.xdg-desktop-portal
